# packaging/cec-remote.service.in
[Unit]
Description=CEC remote bridge and active source handler
After=multi-user.target

[Service]
Type=simple
Environment=HOME=__RUN_HOME__
Environment=CEC_REMOTE_CONFIG=__RUN_HOME__/.config/cec-tv-remote/config.toml
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/cec_remote.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
