# packaging/cec-remote.service.in
[Unit]
Description=CEC remote bridge and active source handler
After=multi-user.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/cec_remote.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
